# deploy/navlys-masternav.service
[Unit]
# NAVLYS — MasterNav (bot Telegram, OPTIONNEL : tu peux piloter par WhatsApp à la place)
Description=NAVLYS MasterNav (bot Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=navlys
WorkingDirectory=/opt/navlys-core
EnvironmentFile=/opt/navlys-core/.env
ExecStart=/opt/navlys-core/venv/bin/python -m navlys_core.masternav
Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
